# site/profile/templates/tomcat/tomcat.service.erb
[Unit]
Description=Apache Tomcat Web Application Container
After=network.target
[Service]
Type=forking
ExecStart=/usr/libexec/tomcat/bin/startup.sh
ExecStop=/usr/libexec/tomcat/bin/shutdown.sh
User=tomcat
Group=tomcat
[Install]
WantedBy=multi-user.target
