# site/profile/templates/tomcat/tomcat.service.erb
[Unit]
Description=Apache Tomcat Web Application Container for <%= @name %> application
After=network.target
[Service]
Type=forking
ExecStart=/var/lib/tomcat/<%= @folder %>/bin/startup.sh
ExecStop=/var/lib/tomcat/<%= @folder %>/bin/shutdown.sh
User=tomcat
Group=tomcat
[Install]
WantedBy=multi-user.target
